# systemd/voting-kiosk.service
[Unit]
Description=School Voting System — Chromium kiosk browser
After=graphical.target voting-backend.service
Requires=voting-backend.service

[Service]
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

# Wait 5 s for the backend to be ready before opening Chromium
ExecStartPre=/bin/sleep 5

ExecStart=/usr/bin/chromium-browser \
  --kiosk \
  --no-sandbox \
  --disable-session-crashed-bubble \
  --disable-infobars \
  --no-first-run \
  --disable-restore-session-state \
  --touch-events=enabled \
  --disable-pinch \
  --overscroll-history-navigation=0 \
  http://127.0.0.1:5000

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voting-kiosk

[Install]
WantedBy=graphical.target
